# PKGS/maintainer-scripts/postinst
#!/usr/bin/env bash
if [ "$1" = "configure" ]; then
    if [ -z $2 ]
    then
        # this is clean install
        previous_version="0.0.0"
    else
        # this is upgrade
        previous_version=$2
    fi
    if dpkg --compare-versions "$previous_version" lt "0.9.1"; then
        echo "Upgrading from version $previous_version to 0.9.1"
        adduser --system --no-create-home --disabled-login rsplayer
        usermod -a -G audio,dialout rsplayer || true
        # optional for hardware access
        usermod -a -G i2c,gpio,spi,audio,input rsplayer || true
        chown -R rsplayer /opt/rsplayer
    fi
    # Ensure network mount base directory exists
    mkdir -p /mnt/rsplayer
    chown rsplayer /mnt/rsplayer
    if dpkg --compare-versions "$previous_version" lt "0.9.7"; then
        echo "Upgrading from version $previous_version to 0.9.7"
        chown rsplayer /opt/rsplayer/env
        chown rsplayer /opt/rsplayer/self.*
    fi
    if dpkg --compare-versions "$previous_version" lt "1.5.0"; then
        echo "Upgrading from version $previous_version to 1.5.0"
        echo "Moving configuration.db to configuration.db.backup_1"
        mv /opt/rsplayer/configuration.db /opt/rsplayer/configuration.db.backup_1
        echo "Go to player Settings page, configure it again and update library to scan dsd and dsf files"
    fi
    if dpkg --compare-versions "$previous_version" lt "2.0.0"; then
        echo "Upgrading from version $previous_version to 2.0.0"
        echo "Database engine changed from sled to fjall. Moving old sled databases aside."
        old_db_dir="/opt/rsplayer/old_sled_dbs_safe_to_delete"
        mkdir -p "$old_db_dir"
        for db in metadata.db queue.db loudness.db play_stats.db; do
            if [ -d "/opt/rsplayer/$db" ]; then
                echo "Moving /opt/rsplayer/$db to $old_db_dir/$db"
                mv "/opt/rsplayer/$db" "$old_db_dir/$db"
            fi
        done
        chown -R rsplayer "$old_db_dir"
        echo "Old sled databases moved to $old_db_dir — delete this directory after verifying the upgrade."
        echo "A full metadata rescan will be needed after first launch."
    fi
    systemctl daemon-reload
    systemctl enable rsplayer
    systemctl start rsplayer
fi
