# PKGS/maintainer-scripts/postinst
#!/usr/bin/env bash
systemctl daemon-reload
systemctl enable rsplayer
systemctl start rsplayer
